# etc/systemd/system/hdparm.service
[Unit]
Description=hdparm sleep

[Service]
Type=oneshot
ExecStart=/home/alkeryn/bin/systemd/hdd

[Install]
WantedBy=multi-user.target
