# etc/systemd/system/start.service
[Unit]
Description=Start Script
[Service]
Type=oneshot
ExecStart=/home/alkeryn/bin/systemd/start
[Install]
WantedBy=multi-user.target

